# conf/isuumo.python.socket
[Unit]
Description=isuumo.python

[Socket]
ListenStream=/run/isuumo/python.sock
User=isucon

[Install]
WantedBy=sockets.target
